# client/supervisord.conf
[supervisord]
nodaemon=true

[program:client]
command=npm start

[program:server]
command=npm run start-server:dev